# systemd/calsync.timer
[Unit]
Wants=calsync.service

[Timer]
OnCalendar=Mon..Fri *-*-* 8..17:0/5:00
Persistent=true

[Install]
WantedBy=timers.target
